# LaunchAgent/com.user.screenshots.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.movescreenshots</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alex/screenshots.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Example: run every 24 hours -->
    <key>StartInterval</key>
    <integer>86400</integer>

    <key>StandardOutPath</key>
    <string>/tmp/move_screenshots.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/move_screenshots.err</string>
</dict>
</plist>
